# hermes-agent/rootfs/etc/s6-overlay/s6-rc.d/hermes/run
#!/command/with-contenv bash
# shellcheck shell=bash
# ==============================================================================
# Start Hermes Agent service (s6-overlay v3 longrun)
# ==============================================================================

HERMES_INSTALL="/opt/hermes-agent"
HERMES_HOME="/data"
HERMES_CONFIG="${HERMES_HOME}/.hermes"
OPTIONS_FILE="/data/options.json"

log_info() { bashio::log.info "$*"; }
log_warning() { bashio::log.warning "$*"; }
log_error() { bashio::log.error "$*"; }

get_config() {
    local key="$1"
    local default="${2:-}"
    jq -r ".$key // \"$default\"" "$OPTIONS_FILE" 2>/dev/null
}

config_exists() {
    local key="$1"
    jq -e ".$key" "$OPTIONS_FILE" >/dev/null 2>&1
}

log_info "Hermes Agent starting..."

# Configure Home Assistant credentials
if config_exists "ha_token"; then
    HA_TOKEN=$(get_config "ha_token")
    if [ -n "$HA_TOKEN" ]; then
        log_info "Configuring Home Assistant..."
        mkdir -p "${HERMES_CONFIG}/secrets"
        echo "{\"api_keys\":{\"homeassistant\":\"$HA_TOKEN\"}}" > "${HERMES_CONFIG}/secrets/credentials.json"
        chmod 600 "${HERMES_CONFIG}/secrets/credentials.json"
        export HASS_TOKEN="$HA_TOKEN"
        log_info "HA token configured."
    fi
else
    log_warning "HA token not configured..."
fi

# Wait for venv (s6-rc may start services before cont-init finishes)
VENV_PATH="${HERMES_INSTALL}/venv"
TIMEOUT=60
ELAPSED=0

while [ ! -d "$VENV_PATH" ]; do
    if [ $ELAPSED -ge $TIMEOUT ]; then
        log_error "Venv not found at ${VENV_PATH} after ${TIMEOUT}s"
        exit 1
    fi
    log_info "Waiting for venv at ${VENV_PATH}... (${ELAPSED}/${TIMEOUT}s)"
    sleep 2
    ELAPSED=$((ELAPSED + 2))
done

cd "${HERMES_INSTALL}"
log_info "Starting Hermes Gateway..."

exec ./venv/bin/hermes gateway run
