# hermes-agent/rootfs/etc/s6-overlay/s6-rc.d/hermes/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Hermes Agent - Home Assistant Addon
# ==============================================================================

HERMES_HOME="/data"
LOG_LEVEL=$(bashio::config 'log_level' 'info')
UPDATE_ON_START=$(bashio::config 'update_on_start' 'false')
GIT_BRANCH=$(bashio::config 'git_branch' 'master')

bashio::log.info "Hermes Agent starting..."

# Configure Home Assistant token if provided
if bashio::config.exists 'ha_token'; then
    HA_TOKEN=$(bashio::config 'ha_token')
    if [ -n "$HA_TOKEN" ]; then
        bashio::log.info "Configuring Home Assistant token..."
        mkdir -p "${HERMES_HOME}/.hermes/secrets"
        echo "{\"api_keys\":{\"homeassistant\":\"$HA_TOKEN\"}}" > "${HERMES_HOME}/.hermes/secrets/credentials.json"
        chmod 600 "${HERMES_HOME}/.hermes/secrets/credentials.json"
        bashio::log.info "HA token configured."
    fi
fi

# Check/update Hermes Agent if enabled
if [ "$UPDATE_ON_START" = "true" ]; then
    bashio::log.info "Checking for Hermes Agent updates..."
    if [ -d "${HERMES_HOME}/hermes-agent/.git" ]; then
        cd "${HERMES_HOME}/hermes-agent"
        git fetch origin "$GIT_BRANCH"
        CURRENT=$(git rev-parse HEAD)
        REMOTE=$(git rev-parse "origin/$GIT_BRANCH")
        if [ "$CURRENT" != "$REMOTE" ]; then
            bashio::log.info "Update found! Pulling latest..."
            git pull origin "$GIT_BRANCH"
            . venv/bin/activate
            pip install -e . --quiet
            bashio::log.info "Hermes Agent updated!"
        else
            bashio::log.info "Hermes Agent already up to date."
        fi
    fi
else
    bashio::log.info "Skipping update check (update_on_start=false)"
fi

# Start Hermes Gateway in a loop (like DuckDNS)
cd "$HERMES_HOME"
bashio::log.info "Starting Hermes Gateway..."

while true; do
    hermes gateway \
        --host 0.0.0.0 \
        --port 8000 \
        --log-level "$LOG_LEVEL" \
        --insecure \
        --hermes-home "$HERMES_HOME"
    
    bashio::log.warning "Hermes Gateway exited! Restarting in 5 seconds..."
    sleep 5
done
