# hermes-agent/rootfs/etc/services.d/hermes-agent/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Start Hermes Agent service
# ==============================================================================

HERMES_HOME="/data"

log_info() { echo "[INFO] $*"; }
log_error() { echo "[ERROR] $*" >&2; }

log_info "Hermes Agent starting..."

# Activate venv
if [ -f "${HERMES_HOME}/hermes-agent/venv/bin/activate" ]; then
    . "${HERMES_HOME}/hermes-agent/venv/bin/activate"
    log_info "Venv activated"
else
    log_error "Venv not found at ${HERMES_HOME}/hermes-agent/venv"
    exit 1
fi

# Write HA token if provided
if [ -f "/data/options.json" ]; then
    HA_TOKEN=$(jq -r '.ha_token // empty' /data/options.json 2>/dev/null)
    if [ -n "$HA_TOKEN" ] && [ "$HA_TOKEN" != "null" ]; then
        log_info "Configuring HA token..."
        mkdir -p "${HERMES_HOME}/.hermes/secrets"
        echo "{\"api_keys\":{\"homeassistant\":\"${HA_TOKEN}\"}}" > "${HERMES_HOME}/.hermes/secrets/credentials.json"
        chmod 600 "${HERMES_HOME}/.hermes/secrets/credentials.json"
        log_info "HA token configured."
    fi
fi

cd "${HERMES_HOME}/hermes-agent"
log_info "Starting Hermes Gateway..."

exec hermes gateway \
    --host 0.0.0.0 \
    --port 8000 \
    --log-level info \
    --insecure \
    --hermes-home "${HERMES_HOME}"
