# hermes-agent/rootfs/etc/services.d/hermes-agent/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Start Hermes Agent service
# ==============================================================================

HERMES_INSTALL="/opt/hermes-agent"
HERMES_HOME="/data"

log_info() { echo "[INFO] $*" >&2; }
log_error() { echo "[ERROR] $*" >&2; }

log_info "Hermes Agent service starting..."

# Wait for venv (s6-rc may start services before cont-init finishes)
VENV_PATH="${HERMES_INSTALL}/venv"
TIMEOUT=60
ELAPSED=0
while [ ! -d "$VENV_PATH" ]; do
    if [ $ELAPSED -ge $TIMEOUT ]; then
        log_error "Venv not found at ${VENV_PATH} after ${TIMEOUT}s"
        exit 1
    fi
    log_info "Waiting for venv at ${VENV_PATH}... (${ELAPSED}/${TIMEOUT}s)"
    sleep 2
    ELAPSED=$((ELAPSED + 2))
done

cd "${HERMES_INSTALL}"
log_info "Starting Hermes Gateway..."

exec ./venv/bin/hermes gateway run
